# scripts/systemd/stella-ai-port-forwards.service
# STELLA Port-Forward Systemd Service
# Automatically manages kubectl port-forwards for production deployment
#
# Installation:
#   1. Edit this file and replace USER and PATH placeholders
#   2. sudo cp stella-ai-port-forwards.service /etc/systemd/system/
#   3. sudo systemctl daemon-reload
#   4. sudo systemctl enable stella-ai-port-forwards
#   5. sudo systemctl start stella-ai-port-forwards
#
# Usage:
#   sudo systemctl status stella-ai-port-forwards
#   sudo systemctl restart stella-ai-port-forwards
#   sudo journalctl -u stella-ai-port-forwards -f

[Unit]
Description=STELLA Kubernetes Port-Forwards
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=/path/to/stella-backend
Environment="HOME=/home/YOUR_USERNAME"
Environment="KUBECONFIG=/home/YOUR_USERNAME/.kube/config"

# Run the monitoring script
ExecStart=/path/to/stella-backend/scripts/monitor-port-forwards.sh

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stella-ai-port-forwards

[Install]
WantedBy=multi-user.target
